# deploy/aware-ai.service
[Unit]
Description=AWARE-AI FastAPI backend
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/aware-ai
EnvironmentFile=/home/ubuntu/aware-ai/.env
ExecStart=/home/ubuntu/aware-ai/.venv/bin/uvicorn test:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
